# mkosi-extra/etc/systemd/system/xiboplayer-kiosk-firstboot.service
[Unit]
Description=Xibo Kiosk First Boot Setup
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/xiboplayer-kiosk-firstboot-done

[Service]
Type=oneshot
ExecStart=/usr/local/bin/xiboplayer-kiosk-firstboot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
