# mkosi-extra/etc/systemd/system/xiboplayer-kiosk-update.timer
[Unit]
Description=Check for xiboplayer-kiosk and player updates

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
